# .config/systemd/user/backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/user/backup.sh
